# edubian/data/rootfs/etc/systemd/system/fstab.service
[Unit]
Description=Fstab service

[Service]
ExecStart=/etc/init.d/fstab.sh

[Install]
WantedBy=multi-user.target
